# bin/scheduler.service
[Unit]
Description=My alarm program
After=network.target

[Service]
ExecStart=/usr/bin/python3 /var/www/sandtable/src/scheduler.py 
WorkingDirectory=/var/www/sandtable
StandardOutput=/var/log/scheduler.log
StandardError=/var/log/scheduler.log
SyslogIdentifier=alarm
User=root
Group=root
Restart=always

[Install]
WantedBy=multi-user.target


